# opt/config/crontab
* * * * * echo "Hello $(date)" >> /opt/cron.log 2>&1
